# src/multinode/lab-cleanup-daemon/dnlab-lab-cleanup.service
[Unit]
Description=dNLab stale lab cleanup reconciler
Documentation=https://github.com/scaci/dnlab-multinode
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/dnlab/src/multinode/venv/bin/dnlab-lab-cleanup start
Restart=always
RestartSec=10
StateDirectory=dnlab-lab-cleanup

[Install]
WantedBy=multi-user.target
